# supervisord.conf
[supervisord]
nodaemon=true

[program:django]
command=/bin/sh -c "/code/docker-entrypoint.sh"

[program:celery]
command=/bin/sh -c "/code/run.sh"